# debian/h3tunnel-server.service
[Unit]
Description=Secure QUIC Tunnel Server
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/tunnel-server
WorkingDirectory=/etc/tunnel
Restart=on-failure

[Install]
WantedBy=multi-user.target
